# deploy/meridian-app.service
# /etc/systemd/system/meridian-app.service
# Runs the Next.js production build.
# Prerequisite: run `npm run build` after each deploy before restarting.

[Unit]
Description=Alara — Next.js App
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/cntp-ops

ExecStart=/usr/bin/node node_modules/.bin/next start -p 3000

EnvironmentFile=/home/deploy/cntp-ops/.env.production
Environment=NODE_ENV=production
Environment=PORT=3000

Restart=always
RestartSec=5s
StartLimitInterval=60s
StartLimitBurst=5

# Allow reasonable memory for Next.js
LimitNOFILE=65536

StandardOutput=journal
StandardError=journal
SyslogIdentifier=alara-app

[Install]
WantedBy=multi-user.target
